# crontab
@reboot			/usr/local/bin/dfrobotups -s -i 20
